# packaging/deb/debian/task-scheduler.postinst
#!/bin/sh
set -e

case "$1" in
  configure)
    systemctl daemon-reload || true
    systemctl reload dbus 2>/dev/null || systemctl restart dbus || true
    systemctl enable --now task-scheduler-daemon.service || true
    update-desktop-database /usr/share/applications 2>/dev/null || true
    gtk-update-icon-cache -f /usr/share/icons/hicolor 2>/dev/null || true
    ;;
esac

#DEBHELPER#
exit 0
